# roles/colossus-legal/templates/colossus-backend.container.j2
# Managed by Ansible — do not edit manually
# App: colossus-legal | Component: backend | Version: {{ deploy_version }}

[Unit]
Description=Colossus Legal Backend ({{ deploy_version }})
After=network-online.target
Wants=network-online.target

[Container]
Image={{ container_registry }}/colossus-backend:{{ deploy_version }}
ContainerName=colossus-backend
PublishPort={{ colossus_legal_backend_port }}:3403
EnvironmentFile={{ app_env_dir }}/colossus-legal-backend.env
Volume={{ colossus_legal_documents_host_path }}:{{ colossus_legal_documents_container_path }}:ro

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target default.target
