# roles/colossus-legal/templates/colossus-frontend.container.j2
# Managed by Ansible — do not edit manually
# App: colossus-legal | Component: frontend | Version: {{ deploy_version }}

[Unit]
Description=Colossus Legal Frontend ({{ deploy_version }})
After=network-online.target colossus-backend.service
Wants=network-online.target

[Container]
Image={{ container_registry }}/colossus-frontend:{{ deploy_version }}
ContainerName=colossus-frontend
PublishPort={{ colossus_legal_frontend_port }}:5473
Volume={{ app_config_dir }}/colossus-legal/config.js:/usr/share/nginx/html/config.js:ro,z
Volume={{ app_config_dir }}/colossus-legal/nginx.conf:/etc/nginx/conf.d/default.conf:ro,z

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target default.target
