# roles/restate/templates/restate-server.container.j2
# Managed by Ansible — do not edit manually
# App: restate | Component: server | Version: {{ deploy_version }}
[Unit]
Description=Restate Durable Execution Server ({{ deploy_version }})
After=network-online.target
Wants=network-online.target

[Container]
Image={{ restate_image_registry }}/{{ restate_image_name }}:{{ deploy_version }}
ContainerName=restate-server

# ── Port publishing ──────────────────────────────────────────
# 8080 — Ingress: apps submit workflows here
# 9070 — Admin API + Web UI
# 5122 — Fabric (node-to-node + Prometheus metrics)
PublishPort={{ restate_ingress_port }}:8080
PublishPort={{ restate_admin_port }}:9070
PublishPort={{ restate_fabric_port }}:5122

# ── Persistent data (Phase A: Podman named volume) ───────────
# Named volume lives at /var/lib/containers/storage/volumes/restate-data/
# on the VM filesystem. Persists across container restarts and VM reboots.
# Followup: convert to virtiofs/ZFS pattern (see session transition).
Volume={{ restate_data_volume_name }}:{{ restate_container_data_path }}:rw

# ── Environment ──────────────────────────────────────────────
# Node identity — CRITICAL: changing this loses data
Environment=RESTATE_NODE_NAME={{ restate_node_name }}
Environment=RESTATE_CLUSTER_NAME={{ restate_cluster_name }}

# Logging
Environment=RESTATE_LOG_FILTER={{ restate_log_filter }}

# Memory pools (see https://docs.restate.dev/server/memory)
# NOTE: Values like "512 MiB" contain a space, so the entire VAR=value
# assignment MUST be wrapped in double quotes — otherwise systemd parses
# only the integer before the space and Restate sees a bare number.
Environment="RESTATE_ROCKSDB_TOTAL_MEMORY_SIZE={{ restate_rocksdb_total_memory_size }}"
Environment="RESTATE_ADMIN__QUERY_ENGINE__MEMORY_SIZE={{ restate_admin_query_engine_memory_size }}"
Environment="RESTATE_BIFROST__RECORD_CACHE_MEMORY_SIZE={{ restate_bifrost_record_cache_memory_size }}"

# Server-wide default retry policy. Equivalent TOML:
#   [invocation.default-retry-policy]
#   initial-interval = "500ms"
#   exponentiation-factor = 2.0
#   max-attempts = 2
#   max-interval = "1m"
#   on-max-attempts = "kill"
# on-max-attempts = "kill" terminates exhausted invocations instead of
# the default "pause" — paused invocations sit forever and confuse
# operators. The per-service admin-API PATCH approach was abandoned
# because it gets wiped on every new service revision.
Environment=RESTATE_DEFAULT_RETRY_POLICY__INITIAL_INTERVAL=500ms
Environment=RESTATE_DEFAULT_RETRY_POLICY__EXPONENTIATION_FACTOR=2.0
Environment=RESTATE_DEFAULT_RETRY_POLICY__MAX_ATTEMPTS=2
Environment=RESTATE_DEFAULT_RETRY_POLICY__MAX_INTERVAL=1m
Environment=RESTATE_DEFAULT_RETRY_POLICY__ON_MAX_ATTEMPTS=kill

# Telemetry opt-out
{% if restate_disable_telemetry %}
Environment=DO_NOT_TRACK=1
{% endif %}

# ── Command-line argument for node name ──────────────────────
# We also pass --node-name as a CLI arg because the docs show this
# pattern and it's the most explicit way to ensure the value is set
# correctly during data directory lookup at startup.
Exec=--node-name={{ restate_node_name }}

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target default.target
