# roles/nvidia-gpu/templates/nvidia-persistence.service.j2
# Managed by Ansible (nvidia-gpu role) — do not edit manually.
#
# Enables NVIDIA persistence mode on every boot so the driver stays resident
# without an X server / display manager. Reboot-surviving: this replaces the
# transient `nvidia-smi -pm 1` header state.
[Unit]
Description=Enable NVIDIA persistence mode
After=systemd-modules-load.service
Wants=systemd-modules-load.service
ConditionPathExists=/dev/nvidia0

[Service]
Type=oneshot
ExecStart=/usr/bin/nvidia-smi -pm 1
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
