# roles/alloy-agent/templates/alloy.container.j2
# /etc/containers/systemd/alloy.container
# Grafana Alloy agent — {{ inventory_hostname }}
# Generated by Ansible alloy-agent role — do not edit manually

[Unit]
Description=Grafana Alloy Agent
After=network-online.target
Wants=network-online.target

[Container]
Image={{ alloy_image }}
ContainerName=alloy

# Run as root inside the container — required for reading /var/log/journal
# (owned by root:systemd-journal on the host). The container already has
# read-only access to /proc, /sys, and / so this doesn't change the
# security posture meaningfully.
User=0

# Disable SELinux label confinement — required for reading /var/log/journal
# on CoreOS (SELinux-enforcing). The container already has read-only access
# to the entire host filesystem, so this doesn't change the security posture.
SecurityLabelDisable=true

# Alloy's HTTP server — Prometheus scrapes metrics here
PublishPort={{ alloy_listen_port }}:{{ alloy_listen_port }}

# Config file (Ansible-managed)
Volume={{ alloy_config_dir_podman }}/config.alloy:/etc/alloy/config.alloy:ro,Z

# Journal access for log shipping
Volume=/var/log/journal:/var/log/journal:ro

# Host filesystem access for metrics collection
# These bind-mounts let the containerized unix exporter read host-level
# /proc, /sys, and filesystem data — same pattern as running node_exporter
# in a container. :ro = read-only, rslave = propagate host mount events.
Volume=/proc:/host/proc:ro
Volume=/sys:/host/sys:ro
Volume=/:/host/root:ro,rslave

# Alloy CLI: "run <config_path> [flags]"
# Note: Alloy uses positional arg for config, not --config.file (lesson from 6A-1 issue #5)
Exec=run /etc/alloy/config.alloy --server.http.listen-addr=0.0.0.0:{{ alloy_listen_port }}

[Install]
WantedBy=multi-user.target
