# roles/vllm-serve/templates/gpu1-profile-restore.service.j2
# Managed by Ansible — do not edit manually
# Boot-restore oneshot: reconcile GPU units to the last-active profile.
# Runs after the Quadlet generator has emitted vllm-*.service and after network.
[Unit]
Description=Restore GPU profile (last-active) on boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart={{ gpu1_profile_restore_script }}

[Install]
WantedBy=multi-user.target
