# roles/vllm-serve/templates/vllm.container.j2
# Managed by Ansible — do not edit manually
# Platform: colossus-llm1 | Layout: {{ item.name }} | Image: {{ vllm_image }}
# Generic GPU-layout unit: serves whatever config sits at the fixed path
# /models/configs/{{ item.name }}.yaml. Model is chosen by cp-ing a per-model
# config onto that path + restarting — the unit itself is model-agnostic.

[Unit]
Description=vLLM {{ item.name }} layout (serves /models/configs/{{ item.name }}.yaml)
After=network-online.target
Wants=network-online.target

[Container]
Image={{ vllm_image }}
ContainerName=vllm-{{ item.name }}
Network=host
{% for g in item.gpus.split(',') %}
AddDevice=nvidia.com/gpu={{ g }}
{% endfor %}
PodmanArgs=--shm-size=8g
EnvironmentFile={{ vllm_env_dir }}/vllm-{{ item.name }}.env
Environment=HF_HOME=/models
Volume={{ vllm_models_dir }}:/models:rw
Exec=vllm serve --config /models/configs/{{ item.name }}.yaml

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
