# roles/traefik/templates/traefik.service.j2
# {{ ansible_managed }}
[Unit]
Description=Traefik Reverse Proxy
Documentation=https://doc.traefik.io/traefik/
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
EnvironmentFile={{ traefik_config_root }}/cloudflare.env
ExecStart={{ traefik_binary_path }} --configFile={{ traefik_config_root }}/traefik.yml
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ traefik_acme_file }} {{ traefik_log_dir }}

[Install]
WantedBy=multi-user.target
